# assets/etc/supervisord/watch_folders.conf
[program:watch_folders]
command = /usr/local/bin/watch_folders
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
